# Assignment04/00-usb-keyboard.rules
ACTION=="add", SUBSYSTEM=="input", ATTRS{configuration}=="HID Keyboard", RUN+="/usr/sbin/insmod /root/Assignment04/main.ko"
